# profile/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty -o '-p -f -- \\u' --noclear --autologin munity - $TERM
ExecStart=/bin/sh -c '/ast/scr/start.sh'
